# dot_config/systemd/user/nm-applet.service
[Unit]
Description=NetworkManager applet

[Service]
ExecStart=/usr/bin/nm-applet
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target