# tools/r2e.service
[Unit]
Description=r2e run
Wants=offlineimap.service
Before=offlineimap.service

[Service]
ExecStart=/usr/bin/r2e run
Type=oneshot