# tools/r2e.service
[Unit]
Description=r2e run

[Service]
ExecStart=/usr/bin/r2e run