# tools/r2e.timer
[Unit]
Description=Launch r2e run every 10 minutes

[Timer]
OnBootSec=3m
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
