# home/dot_config/systemd/user/xwayland-satellite.service
[Unit]
Description=xwayland-satellite
BindsTo=graphical-session.target
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStartPre=/bin/sleep 1
ExecStart=/usr/sbin/xwayland-satellite
Type=notify
NotifyAccess=all
StandardOutput=journal
Restart=always
RestartSec=5s

[Install]
WantedBy=graphical-session.target
